# config/com.gbrain.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gbrain.serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/luca/.bun/bin/bun</string>
        <string>run</string>
        <string>server.ts</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/luca/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/luca</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/luca/Desktop/gbrain-capture</string>

    <key>StandardOutPath</key>
    <string>/Users/luca/Library/Logs/gbrain-capture.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/luca/Library/Logs/gbrain-capture.log</string>
</dict>
</plist>
